# sparse/usr/lib/systemd/system/systemd-ask-password-gui.service.d/20-droid.conf
[Unit]
After=linkerconfig-system.service
Requires=linkerconfig-system.service
