# sparse/usr/lib/systemd/system/usb-moded.service.d/10-wait-for-power.conf
[Service]
# Wait for the /sys/class/power_supply/usb presence by using a file inside the dir.
ExecStartPre=/usr/lib/startup/wait-for-file /sys/class/power_supply/usb/current_now
